# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
directory = /opt/swagip/git
command = /opt/swagip/env/bin/uwsgi --ini=/opt/swagip/git/wsgi.ini
